# packaging/assets/zswap-cli.service.in
[Unit]
Description=ZSwap-cli service
Documentation=man:zswap-cli(1) https://github.com/xvitaly/zswap-cli

[Service]
Type=oneshot
EnvironmentFile=@CMAKE_INSTALL_FULL_SYSCONFDIR@/zswap-cli.conf
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/zswap-cli --env
ExecStop=@CMAKE_INSTALL_FULL_SBINDIR@/zswap-cli -e N
RemainAfterExit=yes
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
