# assets/systemd/zswap-cli.service.in
#
# SPDX-FileCopyrightText: 2020-2025 EasyCoding Team and contributors
#
# SPDX-License-Identifier: MIT
#

[Unit]
Description=@CMAKE_PROJECT_DESCRIPTION@
Documentation=man:@APP_NAME@(1) @CMAKE_PROJECT_HOMEPAGE_URL@

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/@APP_NAME@ --config @CMAKE_INSTALL_FULL_SYSCONFDIR@/@APP_NAME@/@CONF_FILE_NAME@
ExecStop=@CMAKE_INSTALL_FULL_BINDIR@/@APP_NAME@ --enable N
RemainAfterExit=yes
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
