# assets/zswap-cli.service.in
[Unit]
Description=@CMAKE_PROJECT_DESCRIPTION@
Documentation=man:@APP_NAME@(1) @CMAKE_PROJECT_HOMEPAGE_URL@

[Service]
Type=oneshot
EnvironmentFile=@CMAKE_INSTALL_FULL_SYSCONFDIR@/@APP_NAME@.conf
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/@APP_NAME@ --env
ExecStop=@CMAKE_INSTALL_FULL_SBINDIR@/@APP_NAME@ -e N
RemainAfterExit=yes
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
